# deploy/systemd/proxy-gateway-tgbot.service
[Unit]
Description=Proxy Gateway Telegram control bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=@BASE_DIR@/runtime/tgbot.env
ExecStart=/usr/bin/python3 @BASE_DIR@/tgbot.py
Restart=on-failure
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
