# deploy/systemd/proxy-gateway-ios-profile@.service
[Unit]
Description=Proxy Gateway iOS profile responder (per-connection)

[Service]
Type=simple
ExecStart=/usr/bin/python3 /root/5GPN/ios-http.py
Environment=WWW_DIR=/root/5GPN/www
StandardInput=socket
StandardOutput=socket
StandardError=journal
User=root
